# iguana/tools/recharge_iguana.cron
*/10 * * * * root /opt/kmdplatform/iguana/tools/recharge.py >/dev/null 2>&1 

